# nginx_advanced/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start nginx service
# ==============================================================================

set -e

bashio::log.info "Merging options & variables for template"
JSON_CONF=$(jq --arg port $(bashio::core.port) \
    '({options: .}) + ({variables: {port: $port}})' \
    /data/options.json)
bashio::log.info "Generating nginx.conf from template in /etc/nginx/nginx.conf.gtpl"
echo $JSON_CONF | tempio \
    -template /etc/nginx/nginx.conf.gtpl \
    -out /etc/nginx/nginx.conf

DHPARAMS_PATH=/data/dhparams.pem

# Generate dhparams
if ! bashio::fs.file_exists "${DHPARAMS_PATH}"; then
    bashio::log.info  "Generating dhparams (this will take some time)..."
    openssl dhparam -dsaparam -out "$DHPARAMS_PATH" 4096 > /dev/null
fi

# start server
bashio::log.info "Running nginx..."
exec nginx -c /etc/nginx/nginx.conf < /dev/null
